# overlay/system/etc/systemd/system/lxc-android-config.service.d/override.conf
[Service]
# FIXME: broken vendor kernel
# Set short firmware loading timeout so container startup doesn't block itself
ExecStartPre=-/bin/bash -c "echo 2 > /sys/class/firmware/timeout"
